# brain/deploy/com.eddie.brain-dream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.eddie.brain-dream</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <!-- Replace with the absolute path to brain/dist/cli.js on the Mini, e.g. /Users/max/mind-lint/brain/dist/cli.js -->
    <string>/REPLACE/mind-lint/brain/dist/cli.js</string>
    <string>dream</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- Set ANTHROPIC_API_KEY here on the Mini by the operator. NEVER commit the real value. -->
    <key>ANTHROPIC_API_KEY</key><string>REPLACE_WITH_KEY_ON_MINI</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>3</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <!-- Replace with the absolute path to the Mini's Logs dir, e.g. /Users/max/Library/Logs/brain-dream.log -->
  <string>/REPLACE/Library/Logs/brain-dream.log</string>
  <key>StandardErrorPath</key>
  <string>/REPLACE/Library/Logs/brain-dream.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
